# dist/supervisor.hubot.conf
[program:hubot]
command=/home/hubot/eveohbot/bin/hubot
directory=/home/hubot/eveohbot
user=hubot
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/hubot.out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stderr_logfile=/var/log/hubot.err.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
stopsignal=TERM
environment=PATH='/usr/bin:/usr/local/bin:/home/hubot/eveohbot/node_modules:/home/hubot/eveohbot/node_modules/hubot/node_modules:/home/hubot/local/bin'
